# system/etc/local.d/00-archlinux.stop
#!/bin/sh

CHROOT="/var/chroot/archlinux"
SHELL="/bin/su lemniskett"

umount -l $CHROOT/boot
umount -l $CHROOT/lib/modules
umount -l $CHROOT/var/lib/dbus
umount -l $CHROOT/tmp
umount -l $CHROOT/proc
umount -l $CHROOT/sys
umount -l $CHROOT/dev
umount -l $CHROOT/run
umount -l $CHROOT/home/lemniskett

