# system/etc/local.d/00-archlinux.start
#!/bin/sh

CHROOT="/var/chroot/archlinux"
SHELL="/bin/su lemniskett"

mount -R /home/lemniskett $CHROOT/home/lemniskett
mount -t proc /proc $CHROOT/proc
mount -R /tmp $CHROOT/tmp
mount -R /sys $CHROOT/sys
mount -R /dev $CHROOT/dev
mount -R /run $CHROOT/run
mount -R /var/lib/dbus $CHROOT/var/lib/dbus

exit 0
